# config/systemd-user/rc.service
[Unit]
Description=Execute ~/.rc at startup

[Service]
Type=oneshot
ExecStart=sh -c "~/.rc"
RemainAfterExit=yes

[Install]
WantedBy=default.target
