# config/systemd-user/rc.service
[Unit]
Description=Execute ~/.rc at startup

[Service]
ExecStart=bash -c "~/.rc"

[Install]
WantedBy=default.target
