# config/profile.d/init-project.sh
# new header
alias nh='bash ~/git/myetc/config/init-project/nh.sh'
# new cpp project
alias ncpp='bash ~/git/myetc/config/init-project/ncpp.sh'
# new c project
alias nc='bash ~/git/myetc/config/init-project/nc.sh'
