# config/profile.d/usr-local-bin.sh
case ":$PATH:" in
  *:/usr/local/bin:*) :;;
  *) export PATH=/usr/local/bin:$PATH;;
esac

